# systemd/trading-dashboard.service
[Unit]
Description=AI Trading Agent Dashboard (React + FastAPI)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=%i
WorkingDirectory=/home/%i/AI_Stock_portfolio
Environment=PYTHONUNBUFFERED=1
Environment=DASHBOARD_HOST=0.0.0.0
Environment=DASHBOARD_PORT=8000
ExecStart=/home/%i/AI_Stock_portfolio/.venv/bin/python scripts/run_dashboard.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
