# systemd/trading-agent.service
[Unit]
Description=AI Scalping Trading Agent (Alpaca Paper)
After=network-online.target ollama.service
Wants=network-online.target

[Service]
Type=simple
User=%i
WorkingDirectory=/home/%i/AI_Stock_portfolio
Environment=PYTHONUNBUFFERED=1
Environment=TZ=America/Chicago
ExecStart=/home/%i/AI_Stock_portfolio/.venv/bin/python -m src.main
Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target
